# code-sandboxing/gunicorn.conf
[program:gunicorn]
command=python3 -m gunicorn.app.wsgiapp codesandbox:app -w 4
directory=/path/to/project
user=<YOUR_USER_NAME>
autostart=true
autorestart=true
redirect_stderr=true
