# Update_Countrys/usr/bin/countryegipt
#!/bin/bash
 
### BEGIN INIT INFO
# Provides: Block country egipt
# Short-Description:  Block country egipt
# Description:  Block country prevent Target Attacks in sistem.
### END INIT INFO
 
 sudo echo "Block Egypt Country";
 sudo ipset restore -! < /etc/sysconfig/ipset.egipt
 sudo iptables -I FORWARD -m set --match-set blockegypt src -j DROP
 sudo iptables -I INPUT -m set --match-set blockegypt src -j DROP;
