# wipe_linux/usr/bin/WipeBSI
#!/bin/bash
 
### BEGIN INIT INFO
# Provides: WipeBSI
# Short-Description: wipe all hard disk
# Description: wipe hard disk complete for extreme security no recovery.
### END INIT INFO
 
 sudo echo "Eliminando informacion del disco duro! bsi";

 sudo scrub -p bsi /dev/sda
 sudo scrub -p bsi /dev/sda1
 sudo scrub -p bsi /dev/sda2
 sudo scrub -p bsi /dev/sda3
 sudo scrub -p bsi /dev/sda4
 sudo scrub -p bsi /dev/sda5
 sudo scrub -p bsi /dev/hda
 sudo scrub -p bsi /dev/hda1
 sudo scrub -p bsi /dev/hda2;
