# wipe_linux/usr/bin/WipeXtreme
#!/bin/bash
 
### BEGIN INIT INFO
# Provides: WipeXtreme
# Short-Description: wipe all hard disk
# Description: wipe hard disk complete for extreme security no recovery.
### END INIT INFO
 
 sudo echo "Eliminando informacion del disco duro! xtreme 35 Pass";

 sudo scrub -p gutmann /dev/sda
 sudo scrub -p gutmann /dev/sda1
 sudo scrub -p gutmann /dev/sda2
 sudo scrub -p gutmann /dev/sda3
 sudo scrub -p gutmann /dev/sda4
 sudo scrub -p gutmann /dev/sda5
 sudo scrub -p gutmann /dev/hda
 sudo scrub -p gutmann /dev/hda1
 sudo scrub -p gutmann /dev/hda2;
