# wipe_linux/usr/bin/WipeMET33
#!/bin/bash
 
### BEGIN INIT INFO
# Provides: WipeMET44
# Short-Description: wipe all hard disk
# Description: wipe hard disk complete for extreme security no recovery.
### END INIT INFO
 
 sudo echo "Eliminando informacion del disco duro! 33-random-pass method";

 sudo scrub -p pfitzner33 /dev/sda
 sudo scrub -p pfitzner33 /dev/sda1
 sudo scrub -p pfitzner33 /dev/sda2
 sudo scrub -p pfitzner33 /dev/sda3
 sudo scrub -p pfitzner33 /dev/sda4
 sudo scrub -p pfitzner33 /dev/sda5
 sudo scrub -p pfitzner33 /dev/hda
 sudo scrub -p pfitzner33 /dev/hda1
 sudo scrub -p pfitzner33 /dev/hda2;
